# srvdir/wifi-config.service
# /etc/systemd/system/wifi-config.service
[Unit]
Description=Apply Wi-Fi Configuration
After=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-for-wifi.sh
